# posts/openvpn-zugriff-auf-netzwerk-hinter-einem-client/openvpn-gateway-init.service
[Unit]
Description=Setup OpenVPN Gateway
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/openvpn-gateway-init.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
